# infrastructure/systemd/geoip-worker.service
# Systemd service for Geo-IP Screenshot Worker
# Install to: /etc/systemd/system/geoip-worker.service

[Unit]
Description=Geo-IP Screenshot Worker (BullMQ)
Documentation=https://github.com/yourusername/geo-ip-webserver
After=network.target redis.service
Wants=redis.service

[Service]
Type=simple
User=geoip
Group=geoip
WorkingDirectory=/opt/geo-ip-webserver/packages/workers

# Environment
Environment=NODE_ENV=production
EnvironmentFile=/opt/geo-ip-webserver/.env

# Start command
ExecStart=/usr/bin/node dist/screenshot-worker.js

# Restart policy
Restart=always
RestartSec=10
StartLimitInterval=60
StartLimitBurst=3

# Security hardening
NoNewPrivileges=true
PrivateTmp=true
ProtectSystem=strict
ProtectHome=true
ReadWritePaths=/opt/geo-ip-webserver/packages/workers/data

# Resource limits
LimitNOFILE=65536
LimitNPROC=4096

# Logging
StandardOutput=journal
StandardError=journal
SyslogIdentifier=geoip-worker

[Install]
WantedBy=multi-user.target
